# units/mochad@.service.d/publish-1099.conf
[Service]
# Expose the TCP port (1099)
Environment='DOCKER_PUBLISH_1099=1099:1099'
